# upstart/fms_socket-faye.conf
start on starting fms_socket
stop on stopping fms_socket
